# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
loglevel=info
pidfile=/var/run/supervisord.pid

; ── VNC Server ──────────────────────────────────────────────────────────────
[program:vncserver]
command=su - %(ENV_USER_NAME)s -c "vncserver :1 -geometry %(ENV_RESOLUTION)s -depth %(ENV_VNC_DEPTH)s -localhost no -fg"
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/vnc.log
stderr_logfile=/var/log/supervisor/vnc_err.log
startretries=3

; ── noVNC / websockify ───────────────────────────────────────────────────────
[program:novnc]
command=websockify --web=/usr/share/novnc/ %(ENV_NOVNC_PORT)s localhost:%(ENV_VNC_PORT)s
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc_err.log
startretries=5
; Give VNC a moment to start before noVNC tries to connect
startsecs=8
