# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
loglevel=info
pidfile=/var/run/supervisord.pid

[program:vncserver]
; -nolisten tcp = no raw TCP, -AlwaysShared = multiple viewers ok
; -ZlibLevel 1 = fastest compression tradeoff for low CPU
command=su - %(ENV_USER_NAME)s -c "vncserver :1 -geometry %(ENV_RESOLUTION)s -depth %(ENV_VNC_DEPTH)s -localhost no -fg -ZlibLevel 1 -nolisten tcp"
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/vnc.log
stderr_logfile=/var/log/supervisor/vnc_err.log
startretries=3

[program:novnc]
command=websockify --web=/opt/novnc/ %(ENV_NOVNC_PORT)s localhost:%(ENV_VNC_PORT)s
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc_err.log
startretries=5
startsecs=8
