# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
loglevel=info
pidfile=/var/run/supervisord.pid

[program:vncserver]
command=su - %(ENV_USER_NAME)s -c "vncserver :1 -geometry %(ENV_RESOLUTION)s -depth %(ENV_VNC_DEPTH)s -localhost no -fg"
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/vnc.log
stderr_logfile=/var/log/supervisor/vnc_err.log
startretries=3

[program:novnc]
command=websockify --web=/usr/share/novnc/ %(ENV_NOVNC_PORT)s localhost:%(ENV_VNC_PORT)s
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc_err.log
startretries=5
startsecs=8

[program:cloudflared]
command=/bin/bash -c 'if [ -n "%(ENV_CLOUDFLARE_TUNNEL_TOKEN)s" ]; then cloudflared tunnel --no-autoupdate run --token %(ENV_CLOUDFLARE_TUNNEL_TOKEN)s; else sleep infinity; fi'
autostart=true
autorestart=true
priority=30
stdout_logfile=/var/log/supervisor/cloudflared.log
stderr_logfile=/var/log/supervisor/cloudflared_err.log
startretries=5
